# scheduler_config.cron
#!/bin/bash
# Cron job to run the watcher script periodically
# This file should be placed in /etc/cron.d/ or added to crontab

# Run the watcher every 5 minutes to check for new tasks
*/5 * * * * /usr/bin/python3 /path/to/AI_Employee_Vault/watcher.py >> /var/log/ai_employee.log 2>&1

# Alternative: Daily backup of the vault
0 2 * * * /bin/tar -czf /backup/vault_backup_$(date +\%Y\%m\%d).tar.gz /path/to/AI_Employee_Vault >> /var/log/ai_employee_backup.log 2>&1